# pillpal.service
[Unit]
Description=PillPal Auto Startup (AP + Server + Kiosk)
After=network.target
Wants=network.target

[Service]
Type=idle
User=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority

# Activate your venv and run threader.py in the background
ExecStart=/bin/bash -c '
source /home/pi/venv/bin/activate
/home/pi/setup-ap.sh
nohup python3 /home/pi/threader.py >/home/pi/threader.log 2>&1 &
# wait for the server to respond
until curl -s http://127.0.0.1:5001 >/dev/null; do sleep 1; done
/home/pi/.config/kiosk/start.sh
'

Restart=always
RestartSec=10

[Install]
WantedBy=graphical.target
